# systemd/oknginx-logrotate.timer
[Unit]
Description=Run oknginx-logrotate.service hourly to cap the nginx access-log size

[Timer]
OnCalendar=hourly
Persistent=true
Unit=oknginx-logrotate.service

[Install]
WantedBy=timers.target
